# docker-chmod.conf
start on startup
task
exec chmod 666 /var/run/docker.sock
